# host/com.usbsysmon.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd user agent (the equivalent of usb-sysmon.service on Linux).

  Build the Swift agent and copy it to ~/.local/bin/usb-sysmon, then replace
  CHANGE_ME below with your macOS username:
      cp host/com.usbsysmon.agent.plist ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.usbsysmon.agent.plist
  Unload with:
      launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.usbsysmon.agent.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.usbsysmon.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/CHANGE_ME/.local/bin/usb-sysmon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>3</integer>

    <key>StandardOutPath</key>
    <string>/tmp/usb-sysmon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/usb-sysmon.err.log</string>
</dict>
</plist>
